# volume-rotary-encoder.service

[Unit]
Description=Volume knob monitor

[Service]
ExecStart=/usr/local/bin/volume-rotary-encoder

[Install]
WantedBy=multi-user.target
